# deploy/supervisord/uwsgi.ini
[program:catalog-uwsgi]
directory=/opt/catalog/
command=/opt/virtualenvs/catalog/bin/uwsgi /opt/catalog/deploy/uwsgi/catalog.ini
user=comses
autostart=true
redirect_stderr=true
